# roles/redis/templates/redis.supervisor.conf.j2
[program:redis]
command=/opt/redis/bin/redis-server /etc/redis/redis.conf
user=redis
stdout_logfile=/var/log/redis/redis-out.log
stderr_logfile=/var/log/redis/redis-err.log
autostart=true
